# cfmu/systemd/cfmuvalidate.socket
[Socket]
ListenStream=/run/vfrnav/validator/socket
SocketMode=0666

[Install]
WantedBy=sockets.target
